# crontab
*/1 * * * * env && cd /data && node dist/app.js > ${STDOUT_LOC} 2> ${STDERR_LOC}


